# crontab
#Crontab Code/Befehle

#Auf raspberry pi 
#mit "crontab -e" erreichbar

#bei Dateien beim Start ausführen

@reboot /Pfad/zur/datei/Zugangscodes.py

@reboot //pfad/zur/datei/piPlant_Sensordaten.py

